# deploy/systemd/codex-discord-app-server-bot.service
[Unit]
Description=Codex Discord App Server Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# 使用前请替换为实际部署账号、家目录与项目目录。
User=REPLACE_WITH_USER
Group=REPLACE_WITH_GROUP
WorkingDirectory=REPLACE_WITH_PROJECT_DIR
Environment=HOME=REPLACE_WITH_HOME_DIR
Environment=PYTHONUNBUFFERED=1
# 请确保 PATH 中可找到 uv、codex 以及其它运行时依赖。
Environment=PATH=REPLACE_WITH_PATH
ExecStart=REPLACE_WITH_PROJECT_DIR/.venv/bin/python REPLACE_WITH_PROJECT_DIR/scripts/dev.py
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT
TimeoutStopSec=30
SyslogIdentifier=codex-discord-bot

[Install]
WantedBy=multi-user.target
